# root/etc/services.d/bird/run
#!/usr/bin/with-contenv bash

echo "Starting Bird."
exec \
    s6-setuidgid abc \
    /usr/local/sbin/bird -f -P /var/run/bird/bird.pid
